# deploy/pi-edge-mcp.service
[Unit]
Description=Pi Edge MCP telemetry server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/pi-edge-mcp.env
ExecStart=/opt/pi-edge-mcp/.venv/bin/pi-edge-mcp
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
# The gpiozero tools reach GPIO via lgpio (/dev/gpiochip*). ProtectSystem=strict
# breaks the pin factory, so use =full (protects /usr,/boot,/etc, leaves devices
# reachable) and grant the gpio group. Drop these two lines if you do not use the
# GPIO tools and want stricter sandboxing.
ProtectSystem=full
SupplementaryGroups=gpio

[Install]
WantedBy=multi-user.target
